# isa-gateway.service
[Unit]
Description=ISA Gateway — Semantic Field Communication Server
After=network.target

[Service]
Type=simple
User=zcs
WorkingDirectory=/home/zcs/projects/isa
ExecStart=/usr/bin/python3 /home/zcs/projects/isa/gateway.py --host 0.0.0.0 --port 8766
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=isa-gateway

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/zcs/.hermes/isa
ReadOnlyPaths=/home/zcs/projects/isa

[Install]
WantedBy=default.target
