# isa-agent.service
[Unit]
Description=ISA Agent — 军师 Semantic Field Agent
After=isa-gateway.service
Requires=isa-gateway.service

[Service]
Type=simple
User=zcs
WorkingDirectory=/home/zcs/projects/isa
ExecStart=/usr/bin/python3 /home/zcs/projects/isa/isa.py --id 军师 --channel main --agently-listen --gateway ws://localhost:8766
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=isa-agent

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/zcs/.hermes/isa
ReadOnlyPaths=/home/zcs/projects/isa

[Install]
WantedBy=default.target
